# compose/app.conf
[program:app]
command=docker-compose up
directory=/root/demo-oauth-page/
redirect_stderr=true
autostart=true
autorestart=true
priority=10